# crontab
*/2 * * * * python3 /app/autoLogger.py /app
